# system_files/etc/systemd/system/mitie-caddy.service
[Unit]
Description=Caddy reverse proxy — mitie-server homelab services
After=network-online.target netbird.service podman.socket
Wants=network-online.target
# Only start if the compose file has been deployed (via just deploy-service caddy).
ConditionPathExists=/var/lib/mitie-services/caddy/compose.yaml

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=podman-compose -f /var/lib/mitie-services/caddy/compose.yaml up -d --pull always
ExecStop=podman-compose -f /var/lib/mitie-services/caddy/compose.yaml down

[Install]
WantedBy=multi-user.target
